# site-template/systemd/atftpd.service.in
[Unit]
Description=Advanced TFTP Server 
Documentation=@DOCURL@

After=network.target atftpd.socket
Requires=atftpd.socket


[Service]
#
EnvironmentFile=@INSTALL_LOCATION@/@APPNAME_CONF@
ExecStart=@INSTALL_LOCATION@/sbin/atftpd $OPTIONS
StandardInput=socket
#

[Install]
WantedBy = multi-user.target

